# crontab
*/1 * * * * root cd /ping-monitor && ./pinger.py

